# systemd/mtproxy-monitor-daily.timer
[Unit]
Description=MTProxy Daily Report at 00:05

[Timer]
OnCalendar=*-*-* 00:05:00
Persistent=true

[Install]
WantedBy=timers.target
